# 7.0/trisquel-desktop-common-data/data/etc/xdg/autostart/xscreensaver.desktop
[Desktop Entry]
Name=Xscreensaver
Comment=Screensaver for X
Exec=xscreensaver -no-splash
Terminal=false
Type=Application
X-GNOME-Autostart-Phase=Applications
X-GNOME-AutoRestart=true
NoDisplay=true

